# dotfiles/profile.d/fixdns
#!/bin/bash -e

alias fixdns="sudo killall -HUP mDNSResponder"
