# dotfiles/profile.d/vim
#!/bin/bash

export EDITOR=vim

alias m=mvim
alias v=vim
